# kiosk-xscreensaver.timer
[Unit]
Description=simulate user activity for xscreensaver

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target
